# ccm_profile/files/profile.d/enable-ruby.sh
#!/bin/bash

# Enable Ruby 2.4.0
source scl_source enable rh-ruby24
